# scan_pi/greenbone/gvm-autoscan.service
[Unit]
Description=Greenbone autoscan (wait ready, then Full and fast)
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/gvm-autoscan.sh

[Install]
WantedBy=multi-user.target